# crontab
MAILTO=mcguire@crsr.net
PATH=$PATH:/home/mcguire/nasaimages/nasaimages-bot
23 0-23/2 * * * (cd /home/mcguire/nasaimages/nasaimages-bot; npm run start)
